# src/svxlink/scripts/svxlink.service.in
[Unit]
Description=SvxLink ham radio service
After=network.target

[Service]
Type=simple
EnvironmentFile=@CMAKE_INSTALL_PREFIX@/default/svxlink
PIDFile=${PIDFILE}
ExecStartPre=-/bin/touch ${LOGFILE}
ExecStartPre=-/bin/chown ${RUNASUSER} ${LOGFILE}
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/svxlink --logfile=${LOGFILE} --config=${CFGFILE} --pidfile=${PIDFILE} --runasuser=${RUNASUSER}
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=60
TimeoutStopSec=10
LimitCORE=infinity
WorkingDirectory=@CMAKE_INSTALL_PREFIX@/svxlink

[Install]
WantedBy=multi-user.target
